# package/android-tools/adbd.service
[Unit]
Description=ADB Daemon

[Service]
Restart=always
ExecStart=adbd

[Install]
WantedBy=multi-user.target