# package/petitboot/fs-overlay/etc/profile.d/ps3.sh
echo "  _____  _____  _____ "
echo " |  _  ||  ___||___  |"
echo " |   __||___  ||___  |"
echo " |__|   |_____||_____|"
echo "            L I N U X "
echo ""
